# automake-across-platform-bin/data/vd_host.target
[Unit]
Description=QingCloud Desktop Host Server For Guest
